# configs/supervisor/postfix.conf
[program:postfix]
command=/usr/sbin/postfix -c /etc/postfix start-fg
redirect_stderr=true
autorestart=true
startsecs=10
stopwaitsecs=55